# cosmo-packager/package-configuration/manager/init/workflow.conf
description "workflow-service agent"
 
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

env RACK_ENV=development
env WF_SERVICE_LOGS_PATH=/var/log/cosmo/blueprints
chdir /opt/manager/cosmo-manager-develop/workflow-service/

script
        exec sudo -u root /opt/ruby/bin/ruby -S /opt/ruby/bin/unicorn -l 0.0.0.0:8101
end script
