# packman/package-configuration/riemann/init/riemann.conf.template
description "Riemann"

start on runlevel [2345]
stop on runlevel [016]

kill timeout 60
respawn
respawn limit 10 5

setuid riemann
setgid riemann

exec /usr/bin/riemann /etc/riemann/riemann.config